# renew_ssl_cert.cron
MAILTO=''

0 0 * * *  root  /usr/local/bin/renew_ssl_cert.sh >> /var/log/renew_ssl_cert.log 2>&1
