# renew_ssl_cert.cron
MAILTO=''

0 0 * * *  root  /usr/bin/certbot renew --nginx 2>&1 | xargs echo "$(date):" >> /var/log/renew_ssl_cert.log
